# docker/rootfs/etc/service/cron/run
#!/bin/sh -e

# pipe stderr to stdout and run cron
exec 2>&1
exec crond -f
